# sipXimbot/bin/sipximbot.service.in
[Unit]
Description=sipximbot
Documentation=http://sipxcom.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipximbot start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipximbot stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target
