# sipXsqa/bin/sipxsqa.service.in
[Unit]
Description=sipxsqa
Documentation=http://sipxcom.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipxsqa start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipxsqa stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target
